# dockerfile/kafkamirror/supervisor/mirror.conf
[program:mirror]
command=/usr/bin/start-mirrormaker.sh
autostart=true
autorestart=true